# etc/supervisord.d/mysqld.conf
[program:mysqld]
command=/usr/bin/mysqld_safe
autorestart=true
stdout_logfile=/var/log/supervisor/mysqld.log
redirect_stderr=true
user=mysql
